# roles/vendors/thumbor/templates/systemd/thumbor.socket.j2
[Unit]
Description=Thumbor Socket

[Socket]
SocketUser=thumbor
SocketGroup=www-data
ListenStream=/run/thumbor/thumbor.sock

[Install]
WantedBy=sockets.target
